# infra/production/systemd/worklenz-backup.timer
[Unit]
Description=Run Worklenz backup nightly

[Timer]
OnCalendar=*-*-* 07:15:00 UTC
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
