# infra/production/systemd/worklenz-backup-monitor.service
[Unit]
Description=Verify Worklenz backup freshness
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/worklenz
ExecStart=/srv/worklenz/scripts/check-backup-age.sh
User=root
