# infra/production/systemd/worklenz-backup.service
[Unit]
Description=Encrypted Worklenz PostgreSQL backup
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/worklenz
ExecStart=/srv/worklenz/scripts/backup.sh daily
User=root
UMask=0077
